# automatic_server_stopper.cron
# This file should be put under /etc/cron.d/.
# Also, the permission should be 644. ($ sudo chmod 0644 /etc/cron.d/automatic_server_stopper.cron)

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

# m h dom mon dow user    command
* * * * *         YOUR_USER_NAME    /YOUR_PYTHON_PATH automatic_server_stopper.py
